# configs/wallpaper.cron
*/30 * * * * /usr/bin/feh --bg-fill --randomize /home/adi/Nextcloud/Wallpapers/*
